# launchd/com.photobrain.process.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.photobrain.process</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Volumes/PhotoVault/web/backend/scripts/process_new_photos.py</string>
    </array>
    <key>StartInterval</key><integer>900</integer>
</dict>
</plist>
